# inc/profile.d/go
if [ -d /usr/local/go-old/bin ]; then
    export PATH="$PATH:/usr/local/go/bin"
fi
